# instfiles/xrdp-ng.service.in
[Unit]
Description=xrdp-ng daemon
Requires=xrdp-ng-sesman.service
After=syslog.target network.target xrdp-ng-sesman.service

[Service]
Type=forking
PIDFile=${RDS_PID_PATH}/xrdp-ng.pid
EnvironmentFile=/etc/sysconfig/xrdp-ng
ExecStart=${RDS_SBIN_PATH}/xrdp-ng $RDS_OPTIONS
ExecStop=${RDS_SBIN_PATH}/xrdp-ng $RDS_OPTIONS --kill

[Install]
WantedBy=multi-user.target

