# instfiles/xrdp-ng-sesman.service.in
[Unit]
Description=xrdp-ng session manager
After=syslog.target network.target
StopWhenUnneeded=true

[Service]
Type=forking
PIDFile=${RDS_PID_PATH}/xrdp-ng-sesman.pid
EnvironmentFile=/etc/sysconfig/xrdp-ng
ExecStart=${RDS_SBIN_PATH}/xrdp-ng-sesman $SESMAN_OPTIONS
ExecStop=${RDS_SBIN_PATH}/xrdp-ng-sesman $SESMAN_OPTIONS --kill

[Install]
WantedBy=multi-user.target

